# worker/systemd/tinnimate-fb-worker.service
[Unit]
Description=TinniMate Facebook Social Listening Worker
Documentation=https://github.com/tinnimate/tinnimate
After=network.target redis.service
Wants=redis.service

[Service]
Type=simple
User=haichu
Group=haichu
WorkingDirectory=/home/haichu/tinnimate/worker
EnvironmentFile=/home/haichu/tinnimate/worker/.env

# Dùng xvfb-run để headful Chromium hoạt động trên server không có GUI.
# --auto-servernum: chọn display number tự động
# Khi chỉ chạy headless: có thể bỏ xvfb-run, dùng `node dist/index.js` thẳng.
ExecStart=/usr/bin/xvfb-run \
  --auto-servernum \
  --server-args="-screen 0 1280x720x24 -nolisten tcp" \
  /usr/bin/node dist/index.js

# Restart settings
Restart=on-failure
RestartSec=15
StartLimitIntervalSec=120
StartLimitBurst=5

# Resource limits — tránh Chromium OOM
MemoryMax=800M
TasksMax=256

# Logging → systemd journal (xem: journalctl -u tinnimate-fb-worker -f)
StandardOutput=journal
StandardError=journal
SyslogIdentifier=fb-worker

# Security hardening
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
